# github-webhook.service
[Unit]
Description=GitHub Webhook Server
After=network.target

[Service]
Type=simple
User=u123
WorkingDirectory=/home/u123/contentzavod
Environment=PYTHONPATH=/home/u123/contentzavod
Environment=GITHUB_WEBHOOK_SECRET=your-secret-token-here
ExecStart=/home/u123/contentzavod/venv/bin/python /home/u123/contentzavod/webhook_server.py
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-webhook

[Install]
WantedBy=multi-user.target
